# service_files/imaqliq-server.socket
[Unit]
Description=Imaqliq Socket

[Socket]
ListenStream=0.0.0.0:5000
KeepAlive=true

[Install]
WantedBy=sockets.target